# src/main/resources/archetype-resources/src/main/assembly/dist/install/__artifactId__.service
[Unit]
Description=${projectName} Service

[Service]
ExecStart=/opt/dans.knaw.nl/${artifactId}/bin/${artifactId} server /etc/opt/dans.knaw.nl/${artifactId}/config.yml
# Java returns 143 even if the SIGTERM was handled correctly.
SuccessExitStatus=143
Environment="INITIAL_HEAP_SIZE=-Xms64m"
Environment="LC_ALL=en_US.UTF-8"

User=${artifactId}
Group=${artifactId}

[Install]
WantedBy=multi-user.target
